# raspberry/42-latency-ttyUSB0.rules
KERNEL=="ttyUSB0", RUN+="/home/alarm/E46Track/raspberry/latency_fixer.sh"
