# raspberry/bbpi.socket
[Unit]
Description=bbpi socket
PartOf=bbpi.service

[Socket]
ListenStream=0.0.0.0:9999

[Install]
WantedBy=sockets.target
